# supervisord-pythondemo.conf
[program:pythondemo]
command=/start-pythondemo.sh
numprocs=1
autostart=true
autorestart=true
